# files/pairbox_ssh.conf
description "bridge to devbox"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
   socat TCP-LISTEN:22222,fork TCP:`docker port devbox 22`
end script