# files/pairbox_ssh.conf
description "bridge to devbox"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
   socat TCP-LISTEN:{{ssh_port}},fork TCP:`docker port {{name}} 22`
end script